# systemd/mobile443-update.service
[Unit]
Description=Обновление списка разрешенных мобильных ASN и блок-листов Traffic Guard
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mobile443-update.sh
