# systemd/mobile443-apply.service
[Unit]
Description=Применение наборов mobile443 и правил межсетевого экрана из локального кэша
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mobile443-apply-cache.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
