# systemd/mobile443-update.service
[Unit]
Description=Refresh mobile443 mobile ASN allowlist and Traffic Guard blocklists
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mobile443-update.sh
